# tools/systemd/finterm-stub.service
[Unit]
# Localhost stub server for finterm.
# Serves auth/profile/credits/MFA endpoints on http://127.0.0.1:8765.
# Stdlib-only Python; no network code. See tools/local_stub/server.py.
Description=finterm localhost stub server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/fin/finterm/tools/local_stub/server.py
Restart=on-failure
RestartSec=2

# Defensive sandboxing — the stub only needs to read its source dir, write
# to its sqlite users.db under ~/.fincept-localhost, and listen on a port.
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=default.target
